# scripts/com.pliny.secrets-sidecar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pliny.secrets-sidecar</string>

    <key>ProgramArguments</key>
    <array>
        <string>python3</string>
        <string>pliny_secrets_sidecar.py</string>
    </array>

    <!-- install_sidecar.sh rewrites WorkingDirectory at install time. -->
    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <!-- Start at login. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Respawn on crash. -->
    <key>KeepAlive</key>
    <true/>

    <!-- Limit log noise. Sidecar's own scrubbing handles key redaction. -->
    <key>StandardOutPath</key>
    <string>/tmp/pliny_secrets_sidecar.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/pliny_secrets_sidecar.log</string>

    <!-- Throttle restart if the process keeps crashing. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Lower priority so it never preempts user-interactive apps. -->
    <key>Nice</key>
    <integer>1</integer>

    <!-- Process type "Background" gets fair scheduler treatment. -->
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
